# obstaging.pih-emr.org/scripts/crontab
# Run every night at 7am UTC
0 7 * * * /home/jmiranda/scripts/gmdh-nightly-status.sh
