# obstaging.pih-emr.org/scripts/crontab
# Run every night at 12:30pm UTC / 8:30am ET
30 12 * * * /home/jmiranda/scripts/gmdh-nightly-status.sh
